# roles/big-service/files/big-service.conf
start on filesystem
exec /app/big-service/big-service.py
