# roles/panda-service/files/panda-service.conf
start on filesystem
exec /app/panda-service/panda-service.py
